# src/trellobot.service
[Unit]
Description=Trellobot Service

[Service]
Type=Simple
ExecStart=/home/cardinal/botnet/target/release/botnet
